# scripts/systemd/kynance.service
[Unit]
Description="Kynance earnings report"

[Service]
ExecStart=/usr/bin/kynance --report
User=kyle
